# ghostboard/rootfs/etc/services.d/ghostboard/run
#!/usr/bin/with-contenv bash
exec python3 -u /app/get_text.py