# rss-bridge/rootfs/etc/services.d/rss-bridge/run
#!/usr/bin/with-contenv bash
set -e  # Exit immediately on error
echo "Starting Apache2..."
exec httpd -D FOREGROUND