# rss-bridge/rootfs/etc/services.d/rss-bridge/run
#!/usr/bin/with-contenv bash
set -e
echo "Starting Apache2..."
exec apache2ctl -D FOREGROUND
