# rss-bridge/rootfs/etc/services.d/rss-bridge/run
#!/usr/bin/with-contenv bash
exec httpd -D FOREGROUND
